# claude/lspmux/com.lspmux.server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.lspmux.server.plist
  Launchd agent for the lspmux multiplexing server.

  lspmux is a language server protocol multiplexer that lets multiple LSP clients
  (e.g., editor sessions) share a single long-lived language server instance,
  reducing startup overhead and memory usage.

  Installed to: ~/Library/LaunchAgents/com.lspmux.server.plist
  Logs written to: {{LOG_DIR}}/server.log

  To load manually after first install:
    launchctl load ~/Library/LaunchAgents/com.lspmux.server.plist

  To check status:
    launchctl list com.lspmux.server

  Template placeholders (substituted by claude/.sync):
    {{LSPMUX_BIN}}  Full path to the lspmux binary
    {{LOG_DIR}}     Log directory (e.g., ~/Library/Logs/lspmux)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lspmux.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{LSPMUX_BIN}}</string>
        <string>server</string>
    </array>

    <!-- Auto-start when the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the server if it exits unexpectedly -->
    <key>KeepAlive</key>
    <true/>

    <!-- Only run in graphical (Aqua) login sessions, not SSH or headless -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/server.log</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/server.log</string>
</dict>
</plist>
