# chezmoi/dot_config/systemd/user/llama-swap.service
[Unit]
Description=llama-swap router (on-demand llama-server backends, single :9000 front)
PartOf=local-llm.target
ConditionPathExists=%h/local-llm/bin/llama-swap
ConditionPathExists=%h/local-llm/configs/llama-swap.yaml
# Bound the restart loop: at most 3 starts per 120s, then give up (must exceed RestartSec × StartLimitBurst = 90s)
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/local-llm/logs %h/local-llm/cache/slots
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
# Listen on all interfaces so the crabbot microk8s n8n pods can reach the OpenAI-
# compatible endpoint via a Service/EndpointSlice pointing at the node IP. Reachable
# on the home LAN (192.168.50.0/24) and tailnet; never the internet (no port-forward
# or Funnel fronts :9000). See crabbot cluster/manifests/llama-swap-endpoint.yaml.
ExecStart=%h/local-llm/bin/llama-swap --config %h/local-llm/configs/llama-swap.yaml --listen 0.0.0.0:9000
Restart=on-failure
RestartSec=30
# Contain runaway memory in this cgroup instead of triggering a global OOM kill:
# hot group (~6.5G) + one swapped 30B (~21G) + headroom
MemoryMax=30G
# Make this service the kernel OOM killer's first victim, sparing system services
OOMScoreAdjust=1000
# LLM stack dies first under PSI pressure — 2026-07-08 livelock; systemd-oomd
# installed by crabbot cluster/thrash-protection.sh
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=60%
ManagedOOMSwap=kill
Nice=5
StandardOutput=append:%h/local-llm/logs/llama-swap.log
StandardError=append:%h/local-llm/logs/llama-swap.log

[Install]
WantedBy=local-llm.target
