# chezmoi/dot_config/systemd/user/local-llm.target
[Unit]
Description=Local LLM stack (always-on workers + LiteLLM proxy)
Documentation=file://%h/local-llm/README.md
Wants=worker-igpu.service worker-cpu.service litellm.service
After=worker-igpu.service worker-cpu.service litellm.service

[Install]
WantedBy=default.target
