# chezmoi/dot_config/systemd/user/worker-vision.service
[Unit]
Description=Vision worker (Qwen2.5-VL-7B on CPU, separate Zen 5c pin)
After=network.target
ConditionPathExists=%h/local-llm/models/Qwen2.5-VL-7B-Instruct-Q4_K_M.gguf
ConditionPathExists=%h/local-llm/models/mmproj-Qwen2.5-VL-7B-Instruct-Q8_0.gguf

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
# Pin to upper Zen 5c cores so it competes less with worker-cpu (which uses 4-11,16-23)
ExecStart=/usr/bin/taskset -c 8-11,20-23 %h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Qwen2.5-VL-7B-Instruct-Q4_K_M.gguf \
  --mmproj %h/local-llm/models/mmproj-Qwen2.5-VL-7B-Instruct-Q8_0.gguf \
  --host 127.0.0.1 --port 8082 \
  --ctx-size 8192 \
  --parallel 2 --cont-batching \
  --batch-size 1024 --ubatch-size 256 \
  --threads 8 \
  --device none \
  --alias local-vision \
  --log-file %h/local-llm/logs/worker-vision.log
Restart=on-failure
RestartSec=5
Nice=15

# Not enabled by default — start manually: systemctl --user start worker-vision
