# chezmoi/dot_config/systemd/user/worker-igpu.service
[Unit]
Description=Sonnet-tier worker (Qwen3-30B-A3B on iGPU via Vulkan)
PartOf=local-llm.target
ConditionPathExists=%h/local-llm/models/Qwen3-30B-A3B-Instruct-2507-Q4_K_M.gguf

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/local-llm/logs %h/local-llm/cache/slots
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
ExecStart=%h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Qwen3-30B-A3B-Instruct-2507-Q4_K_M.gguf \
  --host 127.0.0.1 --port 8080 \
  --ctx-size 98304 \
  --rope-scaling yarn --rope-scale 3.0 --yarn-orig-ctx 32768 \
  --parallel 1 --cont-batching \
  --batch-size 4096 --ubatch-size 2048 \
  --cache-type-k q8_0 --cache-type-v q8_0 \
  --n-gpu-layers 99 \
  --flash-attn on \
  --slot-save-path %h/local-llm/cache/slots \
  --alias local-sonnet \
  --log-file %h/local-llm/logs/worker-igpu.log
Restart=on-failure
RestartSec=5
Nice=5

[Install]
WantedBy=local-llm.target
