# chezmoi/dot_config/systemd/user/worker-opus.service
[Unit]
Description=Opus-tier worker (Llama 3.3 70B on iGPU — displaces the llama-swap pool)
# Mutual exclusion for the single iGPU slot: the swap pool (sonnet/coder) lives
# inside llama-swap.service, so the manual 70B must stop the whole router.
Conflicts=llama-swap.service
ConditionPathExists=%h/local-llm/models/Llama-3.3-70B-Instruct-Q4_K_M.gguf
# Bound the restart loop: at most 3 starts per 120s, then give up (must exceed RestartSec × StartLimitBurst = 90s)
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/local-llm/logs %h/local-llm/cache/slots
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
ExecStart=%h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Llama-3.3-70B-Instruct-Q4_K_M.gguf \
  --host 127.0.0.1 --port 8090 \
  --ctx-size 16384 \
  --parallel 1 --cont-batching \
  --batch-size 2048 --ubatch-size 512 \
  --cache-type-k q8_0 --cache-type-v q8_0 \
  --n-gpu-layers 99 \
  --flash-attn on \
  --no-mmap \
  --slot-save-path %h/local-llm/cache/slots \
  --alias local-opus \
  --log-file %h/local-llm/logs/worker-opus.log
Restart=on-failure
RestartSec=30
# Contain runaway memory in this cgroup instead of triggering a global OOM kill
MemoryMax=48G
# Make this worker the kernel OOM killer's first victim, sparing system services
OOMScoreAdjust=1000
Nice=15

# Not enabled by default — start manually: systemctl --user start worker-opus
