# chezmoi/dot_config/systemd/user/worker-opus.service
[Unit]
Description=Opus-tier worker (Llama 3.3 70B on iGPU — displaces Sonnet)
# Mutual exclusion with Sonnet — systemd stops worker-igpu when this starts
Conflicts=worker-igpu.service
ConditionPathExists=%h/local-llm/models/Llama-3.3-70B-Instruct-Q4_K_M.gguf

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/local-llm/logs
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
ExecStart=%h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Llama-3.3-70B-Instruct-Q4_K_M.gguf \
  --host 127.0.0.1 --port 8090 \
  --ctx-size 16384 \
  --parallel 1 --cont-batching \
  --batch-size 1024 --ubatch-size 256 \
  --cache-type-k q8_0 --cache-type-v q8_0 \
  --n-gpu-layers 99 \
  --alias local-opus \
  --log-file %h/local-llm/logs/worker-opus.log
Restart=on-failure
RestartSec=5
Nice=15

# Not enabled by default — start manually: systemctl --user start worker-opus
