# chezmoi/dot_config/systemd/user/worker-cpu.service
[Unit]
Description=Haiku-tier worker (Qwen3-8B on CPU, Zen 5c cores)
PartOf=local-llm.target
ConditionPathExists=%h/local-llm/models/Qwen3-8B-Q4_K_M.gguf
# Bound the restart loop: at most 3 starts per 120s, then give up (must exceed RestartSec × StartLimitBurst = 90s)
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p %h/local-llm/logs %h/local-llm/cache/slots
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
# Pin to Zen 5c efficiency cores (4-11) + their SMT siblings (16-23) — leaves Zen 5 perf cores free for foreground
ExecStart=/usr/bin/taskset -c 4-11,16-23 %h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Qwen3-8B-Q4_K_M.gguf \
  --host 127.0.0.1 --port 8081 \
  --ctx-size 16384 \
  --parallel 1 --cont-batching \
  --batch-size 2048 --ubatch-size 512 \
  --threads 12 \
  --device none \
  --reasoning off \
  --slot-save-path %h/local-llm/cache/slots \
  --alias local-haiku \
  --log-file %h/local-llm/logs/worker-cpu.log
Restart=on-failure
RestartSec=30
# Contain runaway memory in this cgroup instead of triggering a global OOM kill
MemoryMax=8G
# Make this worker the kernel OOM killer's first victim, sparing system services
OOMScoreAdjust=1000
Nice=10

[Install]
WantedBy=local-llm.target
