# chezmoi/dot_config/systemd/user/worker-cpu.service
[Unit]
Description=Haiku-tier worker (Qwen3-8B on CPU, Zen 5c cores)
After=network.target
PartOf=local-llm.target

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=%h/local-llm/bin/llama.cpp"
# Pin to Zen 5c efficiency cores (4-11) + their SMT siblings (16-23) — leaves Zen 5 perf cores free for foreground
ExecStart=/usr/bin/taskset -c 4-11,16-23 %h/local-llm/bin/llama.cpp/llama-server \
  --model %h/local-llm/models/Qwen3-8B-Q4_K_M.gguf \
  --host 127.0.0.1 --port 8081 \
  --ctx-size 16384 \
  --parallel 4 --cont-batching \
  --batch-size 1024 --ubatch-size 256 \
  --threads 12 \
  --device none \
  --reasoning off \
  --alias local-haiku \
  --log-file %h/local-llm/logs/worker-cpu.log
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=local-llm.target
